# smtp2telegram.service
[Unit]
Description=SMTP server bridge to Telegram

[Service]
ExecStart=/usr/bin/smtp2telegram.py

[Install]
WantedBy=multi-user.target
